# src/supervisord.conf
[supervisord]
nodaemon=true
[program:kodi-server]
command=/opt/kodi-server/lib/kodi/kodi.bin --nolirc --headless -p

